# remote/cron.d/dash
#* * * * * dash cd /home/dash/bitcoin-abe && python -m Abe.abe --config abe.conf --no-serve --upgrade
*/5 * * * * dash php /home/dash/restart/check.php >/dev/null 2>&1
*/5 * * * * dash php /home/dash/restart/dashd.php >/dev/null 2>&1
*/2 * * * * dash php /var/www/dash.org.ru/site/private/cron/pay.php >/dev/null 2>&1
0 12 * * * dash php /var/www/dash.org.ru/site/private/cron/price.php >/dev/null 2>&1
*/2 * * * * www-data /usr/bin/curl --silent 'http://dash.org.ru/pages/offline.php' >/dev/null 2>&1
